# gcp/crontab
0 * * * * /home/oscar_tienlong_chung_gmail_com/ws/wind_sensor/gcp/generate_charts.sh

30 * * * * /home/oscar_tienlong_chung_gmail_com/ws/wind_sensor/gcp/run_flask.sh
